# examples/cloudflare_ddns_updater.timer
[Unit]
Description=ddns entry updater for IP %i

[Timer]
OnCalendar=*-*-* *:01:42
# depends on sysdtimer.service

[Install]
WantedBy=timers.target
